# mgboost-panel.service
[Unit]
Description=MGBoost Panel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mgboost-panel
ExecStart=/usr/bin/python3 /opt/mgboost-panel/main.py
EnvironmentFile=/opt/mgboost-panel/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
